# sbbs.service
[Unit]
Description=Run the sbbs backup program

[Service]
Type=simple
ExecStart=/usr/bin/sbbs
